# packages/guardian-heartbeat/guardian-heartbeat.service
[Unit]
Description=Guardian Heartbeat Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/guardian-heartbeat
StandardOutput=journal
StandardError=journal
# Don't fail if heartbeat fails
SuccessExitStatus=0 1 2

[Install]
WantedBy=multi-user.target
